# ops/systemd/grabdiscount-v2.service
[Unit]
Description=GrabDiscount Dashboard v2 (nouvelle archi app/)
After=network.target
# Ne pas démarrer avant le legacy — ils partagent /data
After=grabdiscount.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/grabdiscount
EnvironmentFile=/root/grabdiscount/.env
# Vars dédiées v2 — évite la collision avec `PORT` du legacy dans le .env
# partagé. Force le binding localhost-only (défense en profondeur).
Environment=DASHBOARD_V2_HOST=127.0.0.1
Environment=DASHBOARD_V2_PORT=5002
ExecStart=/root/grabdiscount/.venv/bin/python run_dashboard_v2.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grabdiscount-v2

# Sécurité systemd
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/data /root/grabdiscount
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
